# web/app/conf.d/novnc.conf
[program:novnc]
command=/app/conf.d/novnc.sh
autorestart=true
