# web/app/conf.d/vnc.conf
[program:vnc]
command=/app/conf.d/vnc.sh
autorestart=true
